# out/spec/puppet/modules/aide/templates/aide.cron.erb
<% if scope["::servername"] -%>
# THIS FILE IS MANAGED BY PUPPET - DO NOT MODIFY!!

<% end -%>
#!/bin/bash

/usr/sbin/aide --check | /bin/mail -s "AIDE integrity check run for <%= scope['::hostname'] -%>" root@<%= scope['::fqdn'] -%>
